# deb/dcached_0.1.0_amd64/DEBIAN/postrm
#!/bin/bash

update-rc.d dcached remove

rm -f /usr/local/bin/dcached /etc/init.d/dcached 

