# deb/dcached_reb/dcached_0.1.0_amd64/DEBIAN/postinst
#!/bin/bash

ln -s /opt/dcached/bin/dcached /usr/local/bin/dcached
ln -s /opt/dcached/dcached /etc/init.d/dcached


update-rc.d dcached defaults


